# scripts/podman/clawcore.container.in
# ClawCore gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into clawcore's ~/.config/containers/systemd/
# {{CLAWCORE_HOME}} is replaced at install time.

[Unit]
Description=ClawCore gateway (rootless Podman)

[Container]
Image=clawcore:local
ContainerName=clawcore
UserNS=keep-id
Volume={{CLAWCORE_HOME}}/.clawcore:/home/node/.clawcore
EnvironmentFile={{CLAWCORE_HOME}}/.clawcore/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
